# oracle_worker/supervisord.conf
[supervisord]
user=root
nodaemon=true

[program:oci-cli]
command=