# oracle_sshd/config/supervisord-sshd.conf
[program:sshd]
command=/usr/sbin/sshd -d
autostart=true
autorestart=true
user=root